# templates/com.brainstack.claude-extras.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brainstack.claude-extras</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mustafa.bharmal/Documents/brainstack/.venv/bin/python</string>
        <string>__BRAIN_ROOT__/tools/sync_claude_extras.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>BRAIN_ROOT</key>
        <string>__BRAIN_ROOT__</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__BRAIN_ROOT__/runtime/logs/claude-extras.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__BRAIN_ROOT__/runtime/logs/claude-extras.stderr.log</string>
</dict>
</plist>
