# templates/com.user.agent-dream.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Nightly dream cycle launchd job.

  Install:
    cp com.user.agent-dream.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.user.agent-dream.plist

  Edit before installing:
    1. Replace REPLACE_HOME with your actual home dir (e.g. /Users/yourname)
    2. Replace REPLACE_PYTHON with the full path to a Python 3.10+ interpreter
       (e.g. /opt/homebrew/bin/python3.13). The system /usr/bin/python3 on
       older macOS is 3.9 and will not work.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.agent-dream</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>bash</string>
        <string>-c</string>
        <string>flock -n REPLACE_HOME/.agent/.brain.lock REPLACE_PYTHON REPLACE_HOME/.agent/memory/auto_dream.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>REPLACE_HOME/.agent/dream.log</string>
    <key>StandardErrorPath</key>
    <string>REPLACE_HOME/.agent/dream.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
