# templates/com.user.agent-dream.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Nightly dream cycle launchd job.

  Install:
    cp com.user.agent-dream.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.user.agent-dream.plist

  Edit before installing:
    1. Replace REPLACE_HOME with your actual home dir (e.g. /Users/yourname)
    2. Replace REPLACE_PYTHON with the full path to a Python 3.10+ interpreter
       (e.g. /opt/homebrew/bin/python3.13). The system /usr/bin/python3 on
       older macOS is 3.9 and will not work.

  Locking note: dream_runner.py acquires the brain-wide lock via Python's
  fcntl.flock — no dependency on the GNU `flock(1)` binary which isn't
  bundled with macOS. If the lock is held by sync.sh, the runner exits 0.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.agent-dream</string>

    <key>ProgramArguments</key>
    <array>
        <string>REPLACE_PYTHON</string>
        <string>REPLACE_HOME/.agent/tools/dream_runner.py</string>
        <string>--brain-root</string>
        <string>REPLACE_HOME/.agent</string>
    </array>

    <!-- launchd's minimal PATH excludes /opt/homebrew/bin. Set HOME so
         dream_runner can resolve $HOME-based paths internally. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>REPLACE_HOME</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>REPLACE_HOME/.agent/dream.log</string>
    <key>StandardErrorPath</key>
    <string>REPLACE_HOME/.agent/dream.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
